# package/webface/files/etc/crontabs/root
# syncs time
15	*/2	* * *	/etc/init.d/ntpclient >/dev/null 2>/dev/null

